# etc/systemd/system/cbbc.service
[Unit]
Description=CBB Client Service

[Service]
User=root
ExecStart=/usr/local/bin/cbbc
Restart=always

[Install]
WantedBy=multi-user.target
